# .config/systemd/user/low-battery.timer
[Unit]
Description=low-battery timer

[Timer]
OnActiveSec=30
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
