# .config/systemd/user/low-battery.service
[Unit]
Description=low-battery

[Service]
ExecStart=%h/dotfiles-hyprland/scripts/low-battery.sh
